# configs/systemd/user/backup.service
[Unit]
Description=Backup Task 

[Service]
Type=oneshot
ExecStart=/bin/cron_backup

